# shadowsocks.health.check.repair.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <!-- 名称，要全局唯一 -->
    <key>Label</key>
    <string>shadowsocks.health.check.repair</string>
    <!-- 命令， 第一个为命令，其它为参数-->
    <key>ProgramArguments</key>
    <array>
      <!-- python3 实际路径，路径必须是绝对路径 -->
      <string>/Users/opsforce/.pyenv/versions/venv-3.7.3/bin/python</string>
      <string>/Users/opsforce/Projects/SSHCR/app.py</string>
    </array>
    <!-- 运行时间 -->
    <key>StartCalendarInterval</key>
    <array>
      <dict>
        <key>Minute</key>
        <integer>10</integer>
        <key>Hour</key>
        <integer>09</integer>
      </dict>
      <dict>
        <key>Minute</key>
        <integer>10</integer>
        <key>Hour</key>
        <integer>13</integer>
      </dict>
      <dict>
        <key>Minute</key>
        <integer>10</integer>
        <key>Hour</key>
        <integer>17</integer>
      </dict>
    </array>
    <!-- 运行间隔，与StartCalenderInterval使用其一，单位为秒 -->
    <!-- <key>StartInterval</key>
    <integer>30</integer> -->
    <!-- 标准输入文件 -->
    <!-- <key>StandardInPath</key>
    <string>/Users/opsforce/Projects/applescript/log/run-in-shadowsocks.health.check.repair.log</string> -->
    <!-- 标准输出文件 -->
    <!-- <key>StandardOutPath</key>
    <string>/Users/opsforce/Projects/applescript/log/run-out-shadowsocks.health.check.repair.log</string> -->
    <!-- 标准错误输出文件 -->
    <!-- <key>StandardErrorPath</key>
    <string>/Users/opsforce/Projects/applescript/log/run-err-shadowsocks.health.check.repair.log</string> -->
  </dict>
</plist>